# deploy/kb-gateway.service
[Unit]
Description=Keyflo Learning KB Gateway (MCP streamable-http)
After=network.target docker.service

[Service]
Type=simple
WorkingDirectory=/mnt/blockstorage/business/Keyflo_AI/08_Development/kb-gateway
EnvironmentFile=/mnt/blockstorage/env/global.env
EnvironmentFile=-/mnt/blockstorage/env/kb-gateway.env
Environment=LANGCHAIN_COURSE_REPO=/root/kb-core
ExecStart=/root/.venv-langchain-course/bin/python -m kb_gateway --transport streamable-http
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
